# systemd/mxrc-rt.service
[Unit]
Description=MXRC Real-Time Process
Documentation=https://github.com/your-org/mxrc
After=network.target
Wants=mxrc-monitor.service

[Service]
Type=notify
ExecStart=/usr/local/bin/mxrc-rt
Restart=on-failure
RestartSec=5s
WatchdogSec=10s
TimeoutStartSec=30s
TimeoutStopSec=30s

# Real-time scheduling (requires CAP_SYS_NICE)
CPUSchedulingPolicy=fifo
CPUSchedulingPriority=80

# CPU affinity for RT cores (cores 2-3)
CPUAffinity=2 3

# Memory locking for RT performance
LimitMEMLOCK=infinity

# NUMA policy
NUMAPolicy=bind
NUMAMask=0

# User and group
User=mxrc
Group=mxrc

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/mxrc /var/log/mxrc /tmp/mxrc

# Capabilities
AmbientCapabilities=CAP_SYS_NICE CAP_IPC_LOCK
CapabilityBoundingSet=CAP_SYS_NICE CAP_IPC_LOCK

[Install]
WantedBy=multi-user.target
